# etc/cron.d/geoipupdate
0 16 * * * root /usr/bin/geoipupdate
5 16 * * * root ~/GeoLite2/update-database.sh
